# supervisor-salt.conf
[program:salt-minion]
command=/usr/bin/start-salt-minion
numprocs=1
autostart=true
autorestart=true

